# ansible-schema/roles/web/templates/etc/systemd/system/php5.6-fpm.service.d/override.conf.j2
[Service]
LimitNOFILE=65536
LimitCORE=0
